# configs/gtunnel-server.service
[Unit]
Description=GameTunnel Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/gtunnel-server -addr :4700 -subnet 10.10.0.0/24 -max 10
# 可选参数（取消注释并替换 YOUR_SECRET）:
# ExecStart=/usr/local/bin/gtunnel-server -addr :4700 -subnet 10.10.0.0/24 -max 10 -password YOUR_SECRET -status-addr :4701
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
